# supervisord.conf
[supervisord]
nodaemon=false

[program:pmcd]
command=start

[program:pmwebd]
command=start