# supervisord.conf
[supervisord]
nodaemon=false

[program:pmcd]
command=pmcd start

[program:pmwebd]
command=pmwebd start